# pimanager/static/pimanager/client_files/pimanager.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0,10,20,30,40,50 * * * * root /usr/bin/python3 /home/pimanager/report_status.py >/var/log/pimanager/report_status.log 2>&1
0,10,20,30,40,50 * * * * root /usr/bin/python3 /home/pimanager/action.py >/var/log/pimanager/action.log 2>&1
@reboot root /usr/bin/python3 /home/pimanager/report_status.py >/var/log/pimanager/report_status.log 2>&1
@reboot root /usr/bin/python3 /home/pimanager/action.py >/var/log/pimanager/action.log 2>&1
